# roles/clawdbot/templates/clawdbot-host.service.j2
[Unit]
Description=Clawdbot AI Gateway
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User={{ clawdbot_user }}
Group={{ clawdbot_user }}
WorkingDirectory={{ clawdbot_home }}

# Environment variables
Environment="PATH={{ clawdbot_home }}/.local/bin:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PNPM_HOME={{ clawdbot_home }}/.local/share/pnpm"
Environment="HOME={{ clawdbot_home }}"
Environment="XDG_RUNTIME_DIR=/run/user/{{ clawdbot_uid_value | default('1000') }}"

# DBus session bus
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/{{ clawdbot_uid_value | default('1000') }}/bus"

# Start command
ExecStart={{ clawdbot_home }}/.local/bin/clawdbot gateway

# Restart policy
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths={{ clawdbot_home }}/.clawdbot
ReadWritePaths={{ clawdbot_home }}/.local

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clawdbot

[Install]
WantedBy=multi-user.target
